# launchd/com.ghaida.lead6.autoupdate.catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LEAD6_PROJECT_DIR</key>
		<string>/Users/ghaida/Desktop/lead</string>
		<key>LEAD6_SOURCE_PATHS</key>
		<string>/Users/ghaida/Library/CloudStorage/GoogleDrive-gf.smartas@gmail.com/My Drive/lead6.xlsx</string>
	</dict>
	<key>Label</key>
	<string>com.ghaida.lead6.autoupdate.catchup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Developer/CommandLineTools/usr/bin/python3</string>
		<string>/Users/ghaida/Library/Application Support/lead6/auto_update.py</string>
		<string>--catch-up-missed</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/ghaida/Library/Logs/lead6-auto-update.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/ghaida/Library/Logs/lead6-auto-update.out.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/ghaida/Desktop/lead</string>
</dict>
</plist>
